# copilot-presence@.service
[Unit]
Description=Copilot Presence Agent (%i)
Requires=network-online.target agent-hub.service
After=network-online.target agent-hub.service

[Service]
Type=simple
EnvironmentFile=/etc/copilot.env
WorkingDirectory=%E{REPO_ROOT}/mcp-integration
ExecStart=/usr/bin/env COPILOT_ID=%i HUB_WS_URL=%E{HUB_WS_URL} node %E{REPO_ROOT}/mcp-integration/scripts/copilot-vertex.js
Restart=always
RestartSec=2
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
